# com.whatsyes.sync.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
    <key>Label</key>
    <string>com.whatsyes.sync</string>

    <key>ProgramArguments</key>
    <array>
        <string>/opt/homebrew/bin/python3</string>
        <string>/Users/tomereliel/Documents/Projects/WhatsYes/sync.py</string>
    </array>

    <key>WorkingDirectory</key>
    <string>/Users/tomereliel/Documents/Projects/WhatsYes</string>

    <!-- Run once a day at 06:00 -->
    <key>StartCalendarInterval</key>
    <dict>
        <key>Hour</key>
        <integer>6</integer>
        <key>Minute</key>
        <integer>0</integer>
    </dict>

    <!-- Also run at load (immediately on login) -->
    <key>RunAtLoad</key>
    <true/>

    <!-- Log output -->
    <key>StandardOutPath</key>
    <string>/tmp/whatsyes-sync.log</string>
    <key>StandardErrorPath</key>
    <string>/tmp/whatsyes-sync.log</string>
</dict>
</plist>
